# surfcam.service
[Unit]
Description=SurfCam Weather Pi — surf spot webcam streamer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/ryanpatton/script.sh
Restart=always
RestartSec=10
# After power loss, allow unlimited restart attempts (default burst limit can stop the service).
StartLimitIntervalSec=0
# HLS drain + thread joins can exceed the default 90s stop timeout on slow storage.
TimeoutStopSec=180
User=ryanpatton

# Credentials and runtime paths — fill in API_KEY and SPOT_ID before enabling
Environment="API_KEY=your_api_key_here"
Environment="SPOT_ID=your_spot_id_here"
Environment="LD_LIBRARY_PATH=/usr/local/lib"

[Install]
WantedBy=multi-user.target
